# roles/dotfiles/files/.config/systemd/user/tmux.service
[Unit]
Description=Start tmux server

[Service]
RemainAfterExit=yes
ExecStart=/usr/bin/tmux start-server
ExecStop=/usr/bin/tmux kill-server

[Install]
WantedBy=started.target
