# user/.config/systemd/user/urxvtd.service
[Unit]
Description=Start the urxvt deamon

[Service]
Type=forking
ExecStart=/usr/bin/urxvtd -q -f -o

[Install]
WantedBy=default.target
